# config/systemd/guardian.service
[Unit]
Description=Intelligent DNS Guardian & Self-Healing Service
After=network.target dnsmasq.service unbound.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/dns
ExecStartPre=/bin/bash /home/dns/setup_firewall.sh
ExecStart=/usr/bin/python3 /home/dns/guardian.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
